# ops/systemd/mesh-router-autodeploy.service
[Unit]
Description=Auto-build and deploy mesh-router to private registry and ArgoCD
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=oneshot
User=mesh-router
Group=mesh-router
WorkingDirectory=/srv/mesh
ExecStart=/usr/local/bin/mesh-router-autodeploy.sh
